# aports/device/device-samsung-manta/device-samsung-manta.post-install
#!/bin/sh

sed -i 's/WATCHDOG_DEV=\"\"/WATCHDOG_DEV=\"\/dev\/watchdog\"/' /etc/conf.d/watchdog 2> /tmp/errorsed
rc-update add watchdog

